# workshop/supervisor/loadtest.conf
[program:loadtest]
process_name=myapplication
command=/bin/sh ~/script-session.sh
